# backups.cron
*/30 * * * * /home/ubuntu/backups/mc_backups.sh
